# systemd/daily-roll-checkin-lunch.timer
[Unit]
Description=Check-in do meio-dia (12:00)

[Timer]
Unit=daily-roll-checkin.service
OnCalendar=Mon-Fri 12:00:00
Persistent=false

[Install]
WantedBy=timers.target
