# systemd/daily-roll-checkin-evening.timer
[Unit]
Description=Check-in fim de expediente (18:00)

[Timer]
Unit=daily-roll-checkin.service
OnCalendar=Mon-Fri 18:00:00
Persistent=false

[Install]
WantedBy=timers.target
