# systemd/daily-roll.timer
[Unit]
Description=Run Obsidian Daily Roll every day at 23:55

[Timer]
OnCalendar=*-*-* 23:55:00
Persistent=true

[Install]
WantedBy=timers.target
